# update-motd.d/91-contract-ua-esm-status
#!/bin/sh
stamp="/var/lib/ubuntu-advantage/messages/motd-esm-service-status"

[ ! -r "$stamp" ] || cat "$stamp"
