# update-motd.d/91-contract-ua-esm-status
#!/bin/sh
esm_stamp="/var/lib/ubuntu-advantage/messages/motd-esm-service-status"

[ ! -r "$esm_stamp" ] || cat "$esm_stamp"

auto_attach_stamp="/var/lib/ubuntu-advantage/messages/motd-auto-attach-status"

[ ! -r "$auto_attach_stamp" ] || cat "$auto_attach_stamp"
