# systemd/ua-auto-attach.service
[Unit]
Description=Ubuntu Advantage auto attach
Before=cloud-config.service
WantedBy=cloud-config.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ua auto-attach
TimeoutSec=0

[Install]
WantedBy=multi-user.target
