# systemd/ua-license-check.path
[Unit]
Description=Ubuntu Advantage Trigger to check for GCP Ubuntu Pro licenses

[Path]
PathExists=/var/lib/ubuntu-advantage/marker-license-check
Unit=ua-license-check.service

[Install]
WantedBy=multi-user.target
