# debian/ubuntu-advantage-tools.cron.hourly
#!/bin/sh

set -e

UA="/usr/bin/ubuntu-advantage"
STAMP_DIR="/var/cache/ubuntu-advantage-tools"
STAMP_FILE="$STAMP_DIR/ubuntu-advantage-status.stamp"

[ -x "$UA" ] || exit 0

[ -d "$STAMP_DIR" ] || mkdir -p "$STAMP_DIR"

"$UA" status > "$STAMP_FILE"

